# dj_celery_demo/conf/django_celery_beat.ini
[program:celery_beat]
command=python manage.py celery beat -l INFO
directory=/root/myproject/dj_celery_demo
environment=PATH="/root/myenvs/my_dj_celery/bin"
stdout_logfile=/root/myproject/dj_celery_demo/logs/celery_beat.log
stderr_logfile=/root/myproject/dj_celery_demo/logs/celery_beat.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=997
