# dj_celery_demo/conf/django_celery_run.ini
[program:celery_runserver]
command=python manage.py runserver 0.0.0.0:8000
directory=/root/myproject/dj_celery_demo
environment=PATH="/root/myenvs/my_dj_celery/bin"
stdout_logfile=/root/myproject/dj_celery_demo/logs/runserver.log
stderr_logfile=/root/myproject/dj_celery_demo/logs/runserver.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=995