# conf/django_celery_worker.ini
[program:celery_worker]
command=python manage.py celery worker -l INFO
directory=/root/myproject/dj_celery_demo
environment=PATH="/root/myenvs/my_dj_celery/bin"
stdout_logfile=/root/myproject/dj_celery_demo/logs/celery_worker.log
stderr_logfile=/root/myproject/dj_celery_demo/logs/celery_worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=998
