# conf/django_celery_flower.ini
[program:celery_flower]
command=python manage.py celery flower
directory=/home/wx/projects/django_celery_demo
environment=PATH="/home/wx/pyenvs/dj_celery/bin"
stdout_logfile=/home/wx/projects/django_celery_demo/logs/celery_flower.log
stderr_logfile=/home/wx/projects/django_celery_demo/logs/celery_flower.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=996
