# conf/django_celery_worker.ini
[program:celery_worker]
command=python manage.py celery worker -l INFO
directory=/home/wx/myprojects/django_celery_demo
environment=PATH="/home/wx/pyenvs/dj_celery/bin"
stdout_logfile=/home/wx/myprojects/django_celery_demo/logs/celery_worker.log
stderr_logfile=/home/wx/myprojects/django_celery_demo/logs/celery_worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=998
